# install/dmarc-report-converter.timer
[Unit]
Description=Run DMARC report converter daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=basic.target
